# supervisor_restic.ini
[supervisord]
nodaemon=true

[program:cron]
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=start_cron.sh

[program:restic_backup]
startsecs=0
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=restic_actions.sh backup

[program:restic_cleanup]
startsecs=0
autostart=false
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=restic_actions.sh cleanup
